# oxf-core/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=false

[program:dotnet]
command=dotnet run